# dataplane/standalone/pkthandler/lucius-pkthandler.service
[Unit]
Description=Lucius Packet Handler
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
ExecStart=/usr/bin/pkthandler --alsologtostderr
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target